# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off";

[program:flask_paypal]
command=gunicorn -w 2 -b 0.0.0.0:4567 payment_server:app
