# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=nginx

[program:apache2]
command=/bin/bash -c "gunicorn -w $GUNICORN_WORKER -b $GUNICORN_ADDRESS payment_server:app"
