# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:flask_paypal]
command=/bin/bash -c "gunicorn -w %(GUNICORN_WORKER)s -b %(GUNICORN_ADDRESS)s payment_server:app"
